# shared/sscript/powertop-autotune.service
[Unit]
Description=PowerTOP Auto-Tune (Cazil System)
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/powertop --auto-tune
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
